# systemd/pixelcontroller-slave.service
[Unit]
Description=Slave Daemon for PixelController
Wants=network-online.target
After=network-online.target remote-fs.target

[Service]
User=dietpi
ExecStart=/usr/local/sbin/pixelcontroller-slave
KillSignal=INT

[Install]
WantedBy=multi-user.target
